# systemd/cursor-workers-autoscale.timer
[Unit]
Description=Run Cursor worker autoscaler every 30 seconds

[Timer]
OnBootSec=1min
OnUnitActiveSec=30s
Persistent=true

[Install]
WantedBy=timers.target
