# ec2/systemd/cursor-workers-autoscale.service
[Unit]
Description=Autoscale Cursor self-hosted workers

[Service]
Type=oneshot
ExecStart=/usr/local/bin/cursor-workers-autoscale
