# ec2/systemd/cursor-workers-metrics.service
[Unit]
Description=Publish Cursor worker metrics to CloudWatch

[Service]
Type=oneshot
ExecStart=/usr/local/bin/cursor-workers-publish-metrics
